# examples/schedule/launchd.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!-- ~/Library/LaunchAgents/com.example.databank.github.plist
     Edit the absolute paths and token, then:
       launchctl load ~/Library/LaunchAgents/com.example.databank.github.plist
       launchctl start com.example.databank.github   (one-shot, to verify) -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.example.databank.github</string>
  <key>ProgramArguments</key>
  <array>
    <string>/path/to/databank/bin/databank</string>
    <string>sweep</string>
    <string>github</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>GITHUB_TOKEN</key>
    <string>ghp_xxx</string>
  </dict>
  <key>StartInterval</key>
  <integer>21600</integer>
  <key>StandardOutPath</key>
  <string>/path/to/databank/logs/sweep.log</string>
  <key>StandardErrorPath</key>
  <string>/path/to/databank/logs/sweep.log</string>
</dict>
</plist>
